# ansible/roles/nut/files/nut-driver-watchdog.service
[Unit]
Description=Restart NUT drivers with orphaned sockets (workaround for upstream #2484)
After=nut-driver.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/nut-driver-watchdog
