# ansible/roles/kubeadm/files/iscsi-session-cleanup.service
[Unit]
Description=Clean up stale iSCSI sessions and orphan node DB entries
After=iscsid.service

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/iscsi-session-cleanup
